# collectors/services/coingecko-collector.service
[Unit]
Description=CoinGecko Crypto Price Collector
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
User=prediction-markets
WorkingDirectory=/opt/prediction-markets/collectors
Environment=HOME=/opt/prediction-markets
EnvironmentFile=/etc/prediction-markets/coingecko-collector.env
ExecStart=/opt/prediction-markets/venv/bin/python3 coingecko_collector.py --update
TimeoutStartSec=600
MemoryHigh=128M
MemoryMax=256M
