# collectors/services/data-health-alert.timer
[Unit]
Description=Data health alert timer — periodic health checks

[Timer]
OnCalendar=*-*-* *:00,30:00 UTC
Persistent=true

[Install]
WantedBy=timers.target
