# collectors/services/kalshi-candles-sync.timer
[Unit]
Description=Kalshi candles sync timer — hourly + daily candle collection

[Timer]
OnCalendar=*-*-* *:30:00 UTC
Persistent=true

[Install]
WantedBy=timers.target
